# scripts/trading-system-trigger.service
[Unit]
Description=Trading System Trigger Module (Continuous Service)
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
Restart=always
RestartSec=10
User=root
WorkingDirectory=/opt/trading-system
ExecStart=/opt/trading-system/bin/trading-system -module=trigger
StandardOutput=append:/opt/trading-system/logs/trigger-module.log
StandardError=append:/opt/trading-system/logs/trigger-module.log

# Environment variables
Environment="BROKER_CONFIG_PATH=/opt/trading-system/config/broker-config.json"
Environment="REDIS_ADDR=localhost:6379"
Environment="LOG_LEVEL=INFO"
Environment="TRIGGER_LOG_PATH=/opt/trading-system/logs/trigger-module.log"
Environment="WORKER_POOL_SIZE=5"
Environment="TRIGGER_CHECK_INTERVAL=1ms"
Environment="TRIGGER_HEALTH_CHECK_INTERVAL=1m"

[Install]
WantedBy=multi-user.target


