# scripts/trading-system-read.service
[Unit]
Description=Trading System Read Module
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trading-system
ExecStart=/opt/trading-system/bin/trading-system -module=read
Restart=always
RestartSec=10
StandardOutput=append:/opt/trading-system/logs/read-module.log
StandardError=append:/opt/trading-system/logs/read-module.log

# Environment variables
Environment="GOOGLE_SHEETS_CREDENTIALS_PATH=/opt/trading-system/config/google-credentials.json"
Environment="GOOGLE_SHEET_ID="
Environment="GOOGLE_SHEET_BUY_RANGE=to_buy!B3:J"
Environment="GOOGLE_SHEET_SELL_RANGE=to_sell!B3:J"
Environment="GOOGLE_SHEETS_REFRESH_INTERVAL=1m"
Environment="REDIS_ADDR=localhost:6379"
Environment="LOG_LEVEL=INFO"
Environment="READ_LOG_PATH=/opt/trading-system/logs/read-module.log"

[Install]
WantedBy=multi-user.target

