# scripts/trading-system-trigger.service
[Unit]
Description=Trading System Trigger Module (via Cron)
After=network.target redis.service
Requires=redis.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/trading-system
ExecStart=/opt/trading-system/bin/trading-system -module=trigger
StandardOutput=append:/opt/trading-system/logs/trigger-module.log
StandardError=append:/opt/trading-system/logs/trigger-module.log

# Environment variables
Environment="BROKER_CONFIG_PATH=/opt/trading-system/config/broker-config.json"
Environment="REDIS_ADDR=localhost:6379"
Environment="LOG_LEVEL=INFO"
Environment="TRIGGER_LOG_PATH=/opt/trading-system/logs/trigger-module.log"
Environment="WORKER_POOL_SIZE=5"

[Install]
WantedBy=multi-user.target


